# nextcloud_backup/rootfs/etc/services.d/nextcould_backup/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Community Hass.io Add-ons: Example
# Runs example1 script
# ==============================================================================


bashio::log.info "Starting Node..."

exec /usr/bin/nextcloud_backup.sh
